# rootfs/var/spool/cron/crontabs/root
CRON_TIMEOUT php -q /data/htdocs/engine.php >> /var/log/nginx/torrentmonitor_cron_error.log 2>&1
1-59/2 * * * * /rclonesync.sh 2>&1
